# profile.d/05-ssh-agent
# -*- shell-script -*-

source "${BASHCONF_LIBRARY_DIR}/daemons/ssh-agent"
ssh-agent-start
ssh-agent-cleanup

