# profile.d/10-dbus-launch
# -*- shell-script -*-

source "${BASHCONF_LIBRARY_DIR}/daemons/dbus"
dbus-start
dbus-cleanup
